# conf/etc/supervisor/conf.d/supervisord_redis.conf
[program:redis]
directory=/var/lib/redis
command=/usr/bin/redis-server /etc/redis.conf
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
